# artifacts/relay-daemon/install/eqso-relay@.service
# Unidad de servicio systemd para el demonio de radioenlace eQSO ASORAPA.
#
# Es una plantilla (@): el nombre de instancia es la sala a la que se conecta.
# Ejemplos:
#   systemctl enable --now eqso-relay@CB
#   systemctl enable --now eqso-relay@ASORAPA
#   systemctl enable --now eqso-relay@PRUEBAS
#
# Para cada instancia debe existir el archivo de configuracion:
#   /etc/eqso-relay/<instancia>.json
#
# Instalacion:
#   sudo cp eqso-relay@.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now eqso-relay@CB

[Unit]
Description=eQSO Radioenlace — sala %i
Documentation=https://github.com/daycart/eqso-linux-client
After=network-online.target sound.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=eqso
Group=audio

# La instancia (%i) determina que archivo de config leer.
# Para un servicio no-template (eqso-relay.service sin @) %i queda vacio
# y el daemon usa "CB" como instancia por defecto.
Environment=RELAY_INSTANCE=%i
Environment=NODE_ENV=production

# Ruta al daemon compilado
WorkingDirectory=/opt/eqso-asorapa/artifacts/relay-daemon

# ---------------------------------------------------------------------------
# PRE-ARRANQUE
# ---------------------------------------------------------------------------

# Pre-arranque 0: matar procesos ALSA residuales ANTES del reset USB.
# Critico: si arecord/aplay siguen corriendo cuando se hace el reset USB,
# el kernel los pone en D-state (sueno ininterrumpible) y ya no mueren
# ni con SIGKILL, bloqueando /dev/snd para el siguiente arranque.
# sleep 0.3 permite que el kernel limpie los file descriptors antes de
# que usb-reset.sh desenganche el dispositivo.
ExecStartPre=+/bin/bash -c 'pkill -KILL -x arecord 2>/dev/null; pkill -KILL -x aplay 2>/dev/null; sleep 0.3; true'

# Pre-arranque 1: desvincular CM108, reset IOCTL USB y reconectar.
# El script primero hace unbind (libera D-state si lo habia), luego
# USBDEVFS_RESET electrico, luego bind. Se ejecuta como root (+) para
# acceder a /dev/bus/usb y /sys/bus/usb/drivers.
ExecStartPre=+/bin/sh /opt/eqso-asorapa/artifacts/relay-daemon/install/usb-reset.sh

# Pre-arranque 2: configurar mixer ALSA del CM108.
# Tras el reset USB los niveles vuelven al default del driver (puede ser 0%).
# El script detecta automaticamente el numero de tarjeta CM108.
ExecStartPre=+/bin/sh /opt/eqso-asorapa/artifacts/relay-daemon/install/alsa-setup.sh

# Pre-arranque 3: detectar el dispositivo ttyACM del CM108 y crear /dev/eqso-ptt.
# El numero de dispositivo (ACM0, ACM1...) varia segun el orden de conexion USB.
# Este script busca el ttyACM cuyo USB parent tiene tarjeta de sonido (CM108)
# y crea el symlink /dev/eqso-ptt apuntando al dispositivo correcto.
# El config usa siempre /dev/eqso-ptt, sin depender del numero asignado por el kernel.
ExecStartPre=+/bin/sh /opt/eqso-asorapa/artifacts/relay-daemon/install/detect-ptt.sh

ExecStart=/usr/bin/node --enable-source-maps dist/main.mjs

# ---------------------------------------------------------------------------
# PARADA
# ---------------------------------------------------------------------------
#
# Flujo de parada de systemd:
#   1. ExecStop  →  corre stop-cleanup.sh como root (desvincula USB del bus,
#                   libera D-state de arecord/aplay, los mata con SIGKILL)
#   2. KillSignal (SIGTERM) →  Node recibe la senal; su handler encuentra
#                   arecord ya muerto y sale limpiamente.
#   3. Si Node no sale en TimeoutStopSec → systemd manda SIGKILL al cgroup.
#   4. ExecStopPost → red de seguridad final (repite el unbind por si acaso).
#
# Por que ExecStop y no solo ExecStopPost:
#   ExecStopPost corre DESPUES de que systemd ya mato los procesos, demasiado
#   tarde para rescatar un D-state. ExecStop corre ANTES del SIGTERM, por lo
#   que el unbind USB ya ha despertado a arecord antes de que Node intente
#   cerrarlo.
#
# Por que unbind y no "authorized 0":
#   "authorized 0" corta alimentacion pero el kernel no cancela los URBs
#   pendientes, dejando los procesos en D-state igual. "unbind" desregistra
#   el dispositivo del subsistema USB: el kernel cancela todos los URBs con
#   error, despertando inmediatamente los procesos bloqueados.

ExecStop=+/bin/sh /opt/eqso-asorapa/artifacts/relay-daemon/install/stop-cleanup.sh

ExecStopPost=+/bin/sh /opt/eqso-asorapa/artifacts/relay-daemon/install/stop-cleanup.sh

# Reinicio automatico: siempre, excepto si se para manualmente
Restart=always
RestartSec=5s

# Salida de log a journald (acceso via: journalctl -u eqso-relay@%i -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eqso-relay-%i

# Seguridad basica (solo lo estrictamente necesario)
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true

# Permite acceso a /dev/snd para ALSA
SupplementaryGroups=audio

# Senales de parada:
# SIGTERM → Node recibe la senal, cierra stdin de aplay, espera que aplay
# vacie su buffer DMA (<=1.5s) y luego llama process.exit(0).
# Esto evita el D-state: el D-state ocurre cuando SIGKILL mata aplay en
# mitad de una escritura USB DMA. Con cierre graceful, aplay termina la
# escritura antes de salir y el kernel no queda bloqueado.
# 20s = 2s (unbind) + 5s (Node cierre graceful) + margen
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
